# vm/root/etc/s6-overlay/s6-rc.d/svc-docker-router/run
#!/command/execlineb -P
/usr/bin/docker-router