# vm/root/etc/s6-overlay/s6-rc.d/svc-tailscaled/run
#!/command/execlineb -P
/app/tailscaled --state=/var/lib/tailscale/tailscaled.state --tun=userspace-networking