# vm/root/etc/s6-overlay/s6-rc.d/svc-coredns/run
#!/command/execlineb -P
cd /config/coredns \
/usr/bin/coredns -dns.port=53